# com.squaregolf.connector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.squaregolf.connector</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/squaregolf-connector</string>
        <string>--enable-gspro</string>
        <string>--gspro-ip</string>
        <string>127.0.0.1</string>
        <string>--gspro-port</string>
        <string>921</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/SquareGolf Connector/daemon-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/Library/Logs/SquareGolf Connector/daemon-stderr.log</string>
    
    <key>WorkingDirectory</key>
    <string>/usr/local/bin</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>

